# deploy/systemd/localmesh.service
[Unit]
Description=LocalMesh all-in-one stack (AdGuard + NPM + Dashboard)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/localmesh/docker
EnvironmentFile=-/opt/localmesh/.env.production
ExecStart=/usr/bin/docker compose -f /opt/localmesh/docker/docker-compose.yml up -d --build
ExecStop=/usr/bin/docker compose -f /opt/localmesh/docker/docker-compose.yml down

[Install]
WantedBy=multi-user.target
