# deploy/systemd/localmesh-stack.service
# Superseded by localmesh.service — kept for upgrade compatibility.
# If this unit is active on your system, disable it and enable localmesh.service instead:
#   systemctl disable --now localmesh-stack.service
#   systemctl enable --now localmesh.service
[Unit]
Description=LocalMesh Docker stack (deprecated: use localmesh.service)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/localmesh/docker
ExecStart=/usr/bin/docker compose -f /opt/localmesh/docker/docker-compose.yml up -d --build
ExecStop=/usr/bin/docker compose -f /opt/localmesh/docker/docker-compose.yml down
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
